# sv/xvfb/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars

exec 2>&1
exec Xvfb :0 -screen 0 1024x768x24 -ac +extension RANDR
